# databases/db2/verify.service
[Unit]
Description=Script de verificação

[Service]
ExecStart=/home/romul/verify.sh

[Install]
WantedBy=multi-user.target